# script/miner/stop_miner.service
[Unit]
Description=Stop miner
After=default.target
[Service]
ExecStart=./stop_miner.sh
[Install]
WantedBy=default.target